# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log

[program:backend]
command=sh -c "cd /app/packages/server/src && PYTHONPATH=/app/packages/server/src python -m uvicorn api.server:app --host 0.0.0.0 --port 8096"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.stdout.log
stderr_logfile=/var/log/supervisor/backend.stderr.log
startsecs=5
startretries=3
